# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:ttyd]
command=ttyd -p 7681 -b "/tty" -W bash --login
user=ubuntu
environment=HOME="/home/ubuntu",USER="ubuntu",NVM_DIR="/home/ubuntu/.nvm"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ttyd.log
stderr_logfile=/var/log/supervisor/ttyd.log

[program:api]
command=node /home/ubuntu/api/server.js
user=ubuntu
environment=HOME="/home/ubuntu",USER="ubuntu",NVM_DIR="/home/ubuntu/.nvm",PORT="3000",STATE_FILE="/home/ubuntu/.claude-sessions/state.json",WORKSPACE_ROOT="/home/ubuntu/workspace",GITHUB_USER="edkief",GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",POD_NAME="%(ENV_POD_NAME)s"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
